# roles/netspeed/templates/network-speed.service.j2
[Unit]
Description=Network speed statistics script

[Service]
Type=oneshot
ExecStart={{ bindir }}/network-speed
