# roles/netspeed/templates/network-speed.timer.j2
[Unit]
Description=Network speed monitor timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=multi-user.target
